# deployment/supervisor-queue-worker.conf
; Copy to /etc/supervisor/conf.d/ then run:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start wix-3dstore-worker
;
; Adjust paths if your app lives elsewhere.

[program:wix-3dstore-worker]
command=php /var/www/wix-3dstore/backend/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/wix-3dstore/backend
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/wix-3dstore/backend/storage/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
